# root/etc/services.d/caddy/run
#!/usr/bin/execlineb -P

with-contenv
caddy run --config /etc/caddy/Caddyfile --adapter caddyfile "$CADDY_ARGS"
